# repos/dind/src/etc/supervisor/conf.d/caddy.conf
[program:caddy]
command=/usr/bin/caddy run -config="/etc/caddy/Caddyfile"
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
startretries=3
