# repos/screencast/scripts/supervisord.conf
# /usr/bin/supervisord -n -c /etc/supervisor/supervisord.conf


; supervisor config file

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
user=root
loglevel=debug
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:websockify]
command=websockify --web /usr/share/novnc 26369 localhost:26370
autorestart=true

[program:vnc]
command=Xtigervnc -geometry 1288x804x24 -rfbauth /root/.vnc/passwd -rfbport 26370 -alwaysshared :0
autorestart=true

[program:playwright]
command=node /dist/index.js
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:screencast]
command=node -r esbuild-register dist/index.js
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0


